# deploy/etc-frontend.service
[Unit]
Description=Encanto Trade Center Next.js Frontend
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/etc_frontend
Environment=NODE_ENV=production
Environment=PORT=3000
EnvironmentFile=/home/ubuntu/etc_frontend/.env.local

ExecStart=/usr/bin/npm run start -- --port 3000 --hostname 127.0.0.1

Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
